# roles/repo/templates/prune.service.j2
# {{ ansible_managed }}
[Unit]
Description=Prune restic repo {{ repo }}
ConditionACPower=true

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic.{{ repo }} prune {{ restic_prune.flags|default('') }}
